# roles/nodeapp/templates/nodeapp.service.j2
[Unit]
Description=node server
After=network.target

[Service]
WorkingDirectory={{ nodeapp_home }}/www
ExecStart=/usr/bin/npm start

[Install]
WantedBy=multi-user.target
